# LaunchAgents/disabled/never_existed/com.02luka.cls.rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.02luka.cls.rotate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>$HOME/02luka/g/tools/cls/cls_rotate_logs.zsh</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>$HOME/02luka/logs/cls_rotate.out.log</string>
  <key>StandardErrorPath</key><string>$HOME/02luka/logs/cls_rotate.err.log</string>
</dict>
</plist>
