# g/launchagents/com.02luka.hub.plist.sample
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.02luka.hub</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>$HOME/02luka/tools/hub_server.cjs</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>30</integer>
  <key>StandardOutPath</key><string>$HOME/02luka/g/logs/hub.launchd.out</string>
  <key>StandardErrorPath</key><string>$HOME/02luka/g/logs/hub.launchd.err</string>
  <key>WorkingDirectory</key><string>$HOME/02luka</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>LUKA_REDIS_URL</key><string>redis://127.0.0.1:6379</string>
    <key>LUKA_HOME</key><string>$HOME/02luka</string>
    <key>HUB_PORT</key><string>8787</string>
  </dict>
</dict>
</plist>
