# LaunchAgents/backups/20251117_051850/com.02luka.shell_subscriber.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.02luka.shell_subscriber</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/icmini/02luka/tools/shell_subscriber.zsh</string>
  </array>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/icmini/02luka/logs/shell_subscriber.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/icmini/02luka/logs/shell_subscriber.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
