# LaunchAgents/backups/20251117_051850/com.02luka.gg.nlp-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LC_ALL</key>
		<string>C</string>
		<key>PATH</key>
		<string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>REDISCLI_AUTH</key>
		<string>gggclukaic</string>
		<key>REDIS_URI</key>
		<string>redis://default:gggclukaic@127.0.0.1:6379</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.gg.nlp-bridge</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/icmini/02luka/tools/gg_nlp_bridge.zsh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/gg_nlp_bridge.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/gg_nlp_bridge.stdout.log</string>
	<key>ThrottleInterval</key>
	<integer>30</integer>
</dict>
</plist>
