# LaunchAgents/backups/20251117_051850/com.02luka.autoapprove.rd.plist
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.02luka.autoapprove.rd</string>
  <key>ProgramArguments</key>
  <array><string>/usr/bin/env</string><string>zsh</string><string>-lc</string>
    <string>/Users/icmini/02luka/tools/autoapprove_rd.zsh</string>
  </array>
  <key>WatchPaths</key>
  <array><string>/Users/icmini/02luka/bridge/outbox/RD</string></array>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/icmini/02luka/telemetry/autoapprove.out.log</string>
  <key>StandardErrorPath</key><string>/Users/icmini/02luka/telemetry/autoapprove.err.log</string>
</dict></plist>
