# LaunchAgents/backups/20251117_051850/com.02luka.cls.reflection.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.02luka.cls.reflection.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>$HOME/02luka/tools/cls/cls_reflect.zsh && $HOME/02luka/tools/cls/cls_brain_update.zsh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>0</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>StandardOutPath</key><string>$HOME/02luka/logs/cls_reflection.out.log</string>
  <key>StandardErrorPath</key><string>$HOME/02luka/logs/cls_reflection.err.log</string>
  <key>RunAtLoad</key><true/>
</dict>
</plist>

