# etc/launchagents/com.02luka.prometheus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.prometheus</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/prometheus</string>
		<string>--config.file=/Users/icmini/LocalProjects/02luka_local_g/g/config/prometheus/prometheus.yml</string>
		<string>--storage.tsdb.path=/Users/icmini/02luka/.prometheus-data</string>
		<string>--web.listen-address=127.0.0.1:9090</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/agent/prometheus.err</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/agent/prometheus.out</string>
</dict>
</plist>
