# g/deploy/launchagents/com.02luka.nlp.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LUKA_HOME</key>
		<string>/Users/icmini/LocalProjects/02luka_local_g/g</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.nlp.dispatcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>redis-cli SUBSCRIBE gg:nlp | while read line; do
    if [[ "$line" == message* ]]; then
      # next line is the channel, next is the payload
      read ch; read payload;
      printf "%s" "$payload" | python3 /Users/icmini/LocalProjects/02luka_local_g/g/misc/nlp/nlp_command_dispatcher.py;
    fi
  done</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/agent/nlp_dispatcher.err</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/agent/nlp_dispatcher.log</string>
</dict>
</plist>
