# g/etc/launchagents/com.02luka.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HEALTH_CANDIDATES</key>
		<string>http://127.0.0.1:4100/health http://127.0.0.1:4100</string>
		<key>HOST</key>
		<string>127.0.0.1</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>PORT</key>
		<string>4100</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.02luka.dashboard</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>cd "$HOME/02luka/boss-api" &amp;&amp; npm run start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/icmini/02luka/logs/agent/dashboard.err</string>
	<key>StandardOutPath</key>
	<string>/Users/icmini/02luka/logs/agent/dashboard.out</string>
</dict>
</plist>
