# g/etc/launchagents/com.02luka.pushgateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.02luka.pushgateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/icmini/bin/pushgateway</string>
    <string>--web.listen-address=127.0.0.1:9091</string>
    <string>--persistence.file=/Users/icmini/.local/share/02luka/pushgateway/metrics.db</string>
    <string>--persistence.interval=5m</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/icmini/02luka/logs/agent/pushgateway.out</string>
  <key>StandardErrorPath</key><string>/Users/icmini/02luka/logs/agent/pushgateway.err</string>
</dict></plist>
