# systemd/units/02luka-optimizer.service
[Unit]
Description=02LUKA Nightly Optimizer

[Service]
Type=oneshot
WorkingDirectory=%h/02luka
Environment=NODE_NO_WARNINGS=1
EnvironmentFile=%h/02luka/config/redis.env
ExecStart=/usr/bin/env node knowledge/optimize/nightly_optimizer.cjs --telemetry g/telemetry/latest_rollup.ndjson
StandardOutput=append:%h/02luka/g/telemetry/optimizer.log
StandardError=append:%h/02luka/g/telemetry/optimizer.err
