# systemd/units/02luka-daily-digest.timer
[Unit]
Description=02LUKA Daily Digest Timer

[Timer]
OnCalendar=09:00
Persistent=true

[Install]
WantedBy=timers.target
