# g/legacy_parent/legacy_parent/legacy_parent/fixed_launchagents/com.02luka.task.bus.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.02luka.task.bus.bridge</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HOME/dev/02luka-repo" &amp;&amp; REDIS_URL='redis://:changeme-02luka@127.0.0.1:6379/0' exec python3 g/bridge/task_bus_bridge.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/icmini/Library/Logs/02luka/task_bus_bridge.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/icmini/Library/Logs/02luka/task_bus_bridge.err</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
