# infra/launchd/com.antigravity.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/icmini/02luka/tools/session_save.zsh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/com.antigravity.snapshot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.antigravity.snapshot.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/icmini/02luka</string>
</dict>
</plist>
