# Library/LaunchAgents/com.02luka.auto.commit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.02luka.auto.commit</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/icmini/02luka/tools/auto_commit_work.zsh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer> <!-- Every 30 minutes (1800 seconds) -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>ThrottleInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/icmini/02luka/logs/auto_commit.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/icmini/02luka/logs/auto_commit.stderr.log</string>
</dict>
</plist>

