# Library/LaunchAgents/com.02luka.cls.workflow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.02luka.cls.workflow</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/icmini/dev/02luka-repo/scripts/codex_workflow_assistant.sh</string>
        <string>--scan</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLS_SHELL</key>
        <string>/bin/bash</string>
        <key>CLS_FS_ALLOW</key>
        <string>/Volumes/lukadata:/Volumes/hd2:/Users/icmini/Documents/Projects</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/icmini/dev/02luka-repo</string>
    
    <key>StandardOutPath</key>
    <string>/Volumes/lukadata/CLS/logs/cls_workflow.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/lukadata/CLS/logs/cls_workflow_error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
