# Library/LaunchAgents/com.02luka.auto-commit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.02luka.auto-commit</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/icmini/02luka/tools/auto_commit.zsh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>3600</integer>  <!-- Auto-commit every hour -->
    
    <key>RunAtLoad</key>
    <false/>  <!-- Don't run immediately on boot -->
    
    <key>StandardOutPath</key>
    <string>/Users/icmini/02luka/logs/auto_commit.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/icmini/02luka/logs/auto_commit.log</string>
</dict>
</plist>
