# g/legacy_parent/fixed_launchagents/com.02luka.mcp.fs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.02luka.mcp.fs</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>export FS_ROOT="$HOME/Library/CloudStorage/GoogleDrive-ittipong.c@gmail.com/My Drive/02luka"; exec "$HOME/.venv/mcpfs/bin/python" "$HOME/dev/02luka-repo/g/tools/mcp_fs_server.py"</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/mcp_fs_py.out</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/mcp_fs_py.err</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
