# LaunchAgents/com.02luka.optimizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.02luka.optimizer</string>
    <key>ProgramArguments</key>
      <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/icmini/02luka/knowledge/optimize/nightly_optimizer.cjs</string>
        <string>--telemetry</string>
        <string>g/telemetry/latest_rollup.ndjson</string>
      </array>
    <key>StartCalendarInterval</key>
      <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
    <key>StandardOutPath</key><string>/Users/icmini/02luka/g/telemetry/optimizer.log</string>
    <key>StandardErrorPath</key><string>/Users/icmini/02luka/g/telemetry/optimizer.err</string>
    <key>WorkingDirectory</key><string>/Users/icmini/02luka</string>
    <key>RunAtLoad</key><true/></dict>
</plist>
