# g/launchd/com.02luka.wo_pipeline_guardrail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.02luka.wo_pipeline_guardrail</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/icmini/02luka/g/tools/wo_pipeline/wo_pipeline_guardrail.zsh</string>
  </array>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/Users/icmini/02luka/logs/wo_pipeline/wo_pipeline_guardrail.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/icmini/02luka/logs/wo_pipeline/wo_pipeline_guardrail.err.log</string>
</dict>
</plist>
